# deploy/systemd/accept-gpt.service
[Unit]
Description=AcceptGPT - tự động duyệt thành viên ChatGPT Business
After=network-online.target xvfb.service
Wants=network-online.target
Requires=xvfb.service

[Service]
User=acceptgpt
Group=acceptgpt
WorkingDirectory=/opt/accept-gpt

# Trỏ Chromium vào Xvfb display :99 (khớp với xvfb.service).
Environment=DISPLAY=:99
Environment=NODE_ENV=production
EnvironmentFile=/opt/accept-gpt/.env

# Chạy bản đã build (npm run build -> dist/).
ExecStart=/usr/bin/node dist/main.js

# Tự restart khi process lỗi; systemd tự chạy lại sau khi VPS reboot.
Restart=always
RestartSec=5

# Ghi log vào journald (xem bằng: journalctl -u accept-gpt -f).
StandardOutput=journal
StandardError=journal

# Siết chặt quyền một chút.
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
